# files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

#ebtables -t broute -A BROUTING -p ! ipv6 -j DROP -i eth1
#brctl addif br-lan eth1
#/etc/init.d/odhcpd stop
#/etc/init.d/odhcpd disable

exit 0
